# Custom/Postapplication
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=99

zzz_default() {

	if [ -f "/etc/zzz-default-settings" ]; then
		/etc/zzz-default-settings
                sleep 2
                rm -rf /etc/zzz-default-settings
	else
		echo "No operation"
	fi
    
	if [ -f "/etc/default-setting" ]; then
		/etc/default-setting
                sleep 2
                rm -rf /etc/default-setting
	else
		echo "No operation"
	fi

}

start() {
	zzz_default
}

