# Custom/Postapplication
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=99

zzz_default() {

	if [ -f "/etc/zzz-default-settings" ]; then
		/etc/zzz-default-settings
                rm -rf /etc/zzz-default-settings
	else
		echo "No operation"
	fi
    
	if [ -f "/etc/default-setting" ]; then
		/etc/default-setting
		uci commit network
		uci commit luci
		uci commit dhcp
		uci commit system
		uci commit firewall
		if [ -f /etc/config/ttyd ]; then
		        uci commit ttyd
		fi
                rm -rf /etc/default-setting
	else
		echo "No operation"
	fi
	
	rm -rf /tmp/luci-*cache*
        /etc/init.d/network restart
	/etc/init.d/uhttpd restart
	/etc/init.d/dnsmasq restart
	/etc/init.d/firewall restart
	exit 0

}

start() {
	zzz_default
}
